# system76-notify.timer
[Unit]
Description=Check for system76 updates

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target